# docker/crontab
## min hour day month weekday command
* * * * * /bin/date --rfc-2822 >> /var/log/cron.log
*/30 * * * * bash /home/node/app/exec-bot.sh >> /var/log/cron.log 2>&1